# target/linux/ramips/base-files/etc/hotplug.d/net/00-ralink-wifi-detect
#!/bin/sh

[ "${ACTION}" = "add" ] && [ "${INTERFACE%%[0-9]}" = "ra" ] && {
  /sbin/wifi config
}
